# vp-engine/deployments/ec2/mindbliss-vp-engine.service
[Unit]
Description=Mindbliss Power VP Engine
Documentation=file:/opt/mindbliss/vp-engine/current/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mindbliss
Group=mindbliss
WorkingDirectory=/opt/mindbliss/vp-engine/current
EnvironmentFile=/etc/vp-engine/app.env
ExecStart=/opt/mindbliss/vp-engine/current/vp-engine
Restart=always
RestartSec=5

NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/log/vp-engine
LimitNOFILE=65536

KillSignal=SIGTERM
TimeoutStopSec=45

[Install]
WantedBy=multi-user.target
