# vp-engine/deploy/vp-kb-indexer.service
# systemd unit — vp-kb-indexer (EC2 WORKER)
# Install:
#   sudo cp vp-kb-indexer.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now vp-kb-indexer
[Unit]
Description=VicionPower KB indexer (support RAG: Postgres -> Qdrant)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=vicion
Group=vicion
EnvironmentFile=/etc/vicion/vp-kb-indexer.env
ExecStart=/opt/vicion/bin/vp-kb-indexer
Restart=on-failure
RestartSec=10
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
